# health-check.service
[Unit]
Description=Health Check Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/health-check
ExecStart=/usr/bin/node /home/ubuntu/health-check/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/health-check/.env

[Install]
WantedBy=multi-user.target
